# systemd/.config/systemd/user/mfin.target
[Unit]
Description=Custom systemd stuff

[Install]
Alias=default.target

