# systemd/.config/systemd/user/mpd.socket
[Unit]
Description=Music Player Daemon Sockets

[Socket]
ListenStream=6600

[Install]
WantedBy=sockets.target
