# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=UDiskie Automounter

[Service]
ExecStart=/usr/bin/udiskie -2
ExecStop=/usr/bin/pkill udiskie

[Install]
WantedBy=default.target
